# scripts/ctv-db-validation.service
[Unit]
Description=CTV Production Database Daily Validation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=daseme
Group=ctvapps
WorkingDirectory=/opt/apps/ctv-bookedbiz-db
Environment="PATH=/opt/apps/ctv-bookedbiz-db/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Load environment variables if available
EnvironmentFile=-/etc/ctv-db-validation.env
ExecStart=/opt/apps/ctv-bookedbiz-db/.venv/bin/python3 /opt/apps/ctv-bookedbiz-db/scripts/validate_production_db.py
TimeoutStartSec=600

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/opt/apps/ctv-bookedbiz-db /tmp
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
CapabilityBoundingSet=

