# scripts/tailscale-serve.service
# Tailscale Serve – main + dev on same Pi (one URL, two paths).
# Configures tailscaled to proxy:
#   https://<pi>.tailnet.ts.net/     → http://127.0.0.1:8000 (main)
#   https://<pi>.tailnet.ts.net/dev  → http://127.0.0.1:5100 (dev)
#
# Syntax is for Tailscale 1.52+: use --set-path for subpath, target = port or URL.
# Pre-1.52 used "tailscale serve set <path> <target>" (different order/args).
#
# Install:
#   sudo cp scripts/tailscale-serve.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tailscale-serve && sudo systemctl start tailscale-serve

[Unit]
Description=Tailscale Serve (main + dev paths)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# 1.52+: --yes = non-interactive; target = port (proxies to http://127.0.0.1:PORT)
ExecStart=/usr/bin/tailscale serve --yes 8000
ExecStart=/usr/bin/tailscale serve --yes --set-path=/dev 5100

[Install]
WantedBy=multi-user.target
